# deploy/polymarket-rn1.service
[Unit]
Description=Polymarket copytrader - RN1
After=network-online.target
Wants=network-online.target

[Service]
WorkingDirectory=/opt/polymarket
ExecStart=/opt/polymarket/.venv/bin/python /opt/polymarket/bot.py
Environment=PYTHONUNBUFFERED=1
Environment=TARGET_ADDRESS=0x2005d16a84ceefa912d4e380cd32e7ff827875ea
Environment=TARGET_NAME=RN1
Environment=EXCEL_PATH=/opt/polymarket/data/polymarket_paper_trades.xlsx
# Live trading stays OFF here. To go live, see DEPLOY.md (and keep the key OFF the server until ready).
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
